# install/macos/share_ext/kr.service
[Unit]
Description=Krypton daemon

[Service]
ExecStart=/usr/bin/krd
Restart=on-failure

[Install]
WantedBy=default.target
