# daemon/launch/krssh-agent.service
[Unit]
Description=kryptonite ssh agent
Before=default.target

[Service]
ExecStart=/usr/local/bin/krssh-agent
Restart=always
