# share/kr.service
[Unit]
Description=Kryptonite daemon

[Service]
ExecStart=/usr/bin/krd
Restart=on-failure

[Install]
WantedBy=default.target